# extras/services/rfkill-unblock/run
#!/bin/sh

# Unblock all wireless devices
rfkill unblock all
